# pacman-accel.service
[Unit]
Description=Pacman Accelerator
After=network-online.target

[Service]
ExecStart=/usr/bin/pacman-accel

[Install]
WantedBy=multi-user.target
